# deploy/trackops-sync.service
[Unit]
Description=TrackOps sync server (REST + SSE)
After=network.target
Documentation=https://github.com/anthropics/claude-code

[Service]
Type=simple
User=trackops
Group=trackops
WorkingDirectory=/opt/trackops/server
Environment=NODE_ENV=production
Environment=PORT=8787
Environment=TRACKOPS_DATA_DIR=/opt/trackops/data
# La autenticación es por contraseña de proyecto (ver POST /api/projects/:id/password).
# Si quieres un gate global adicional, define TRACKOPS_TOKEN aquí; si no, déjalo sin
# establecer: las lecturas son públicas y las escrituras requieren contraseña.
# Opcional: añade GEMINI_API_KEY=... en /etc/trackops/env para activar el briefing IA.
EnvironmentFile=-/etc/trackops/env
ExecStart=/usr/bin/node index.js
Restart=on-failure
RestartSec=5

# Endurecido básico
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/trackops/data
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
